# deploy/systemd/agentrunner.service
[Unit]
Description=AgentRunner Discord AI game-development runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/agentrunner
EnvironmentFile=/opt/agentrunner/.env
ExecStart=/usr/local/bin/bun run start
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30
User=agentrunner
Group=agentrunner
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/opt/agentrunner /opt/game-projects /opt/obsidian-vaults

[Install]
WantedBy=multi-user.target
